# libreoffice/files/supervisor.conf
[supervisord]
nodaemon=true

[program:soffice]
command=/usr/bin/soffice --nologo --norestore --invisible --headless --accept='socket,host=0.0.0.0,port=8999,tcpNoDelay=1;urp;'
